# python-anti-cheat/supervisor.conf
; Supervisor configuration for the LearnAdapt Anti-Cheat Python service
; Place this file at /etc/supervisor/conf.d/learnadapt-anticheat.conf
; then run:  sudo supervisorctl reread && sudo supervisorctl update

[program:learnadapt-anticheat]
command=/path/to/learnadapt/python-anti-cheat/venv/bin/uvicorn main:app --host 127.0.0.1 --port 8765 --workers 1
directory=/path/to/learnadapt/python-anti-cheat
user=www-data
autostart=true
autorestart=true
startretries=5
startsecs=5
redirect_stderr=true
stdout_logfile=/var/log/learnadapt-anticheat.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
environment=PYTHONUNBUFFERED="1"
